# start.conf
[program:driver-simulation]
command=php /var/www/artisan s:s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
